# services/reminders/supervisor.conf
[program:reminders]
command=/usr/bin/env python services/reminders/reminders.py %(ENV_FLAGS)s
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
stdout_logfile=log/reminders.log
redirect_stderr=true
