# services/fink/supervisor.conf
[program:fink]
command=/usr/bin/env python services/fink/fink.py %(ENV_FLAGS)s
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
stdout_logfile=log/fink.log
redirect_stderr=true
