# services/mcp_server/supervisor.conf
[program:mcp_server]
command=/usr/bin/env python -m services.mcp_server
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1",MCP_TRANSPORT="http"
stdout_logfile=log/mcp_server.log
redirect_stderr=true
